# Modul1/praktikum-2025/task-3/crontabs
*/2 * * * * /bin/bash /home/mouis/work/praktikum-modul1-d02/task-3/automation.sh
