# Modul1/praktikum-2025/task-4/crontabs
*/5 * * * * /bin/bash /home/mouis/work/praktikum-modul-1-d02/1-d02/task-4/minute5_log.sh

0 * * * * /bin/bash /home/mouis/work/praktikum-modul-1-d02/td02/task-4/agg_5min_to_hour.sh

0 * * * * /bin/bash /home/mouis/work/praktikum-modul-1-d02/td02/task-4/uptime_monitor.sh

0 0 * * * /bin/bash /home/mouis/work/praktikum-modul-1-d02/td02/task-4/cleanup_log.sh
